# source/lumus/etc/systemd/system/at-boot-lumus.service
[Unit]
After=display-manager.service

[Service]
ExecStart=/usr/bin/at-boot-lumus

[Install]
WantedBy=default.target

