# source/lumus/etc/systemd/system/at-boot.service
[Unit]
After=display-manager.service

[Service]
ExecStart=/usr/bin/at-boot

[Install]
WantedBy=default.target

